# components/media/minidlna/Solaris/minidlna.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<service_bundle type='manifest' name='minidlna'>

  <service
      name='network/minidlna'
      type='service'
      version='1'>

    <!--
	Because we may have multiple instances of network/http
	provided by different implementations, we keep dependencies
	and methods within the instance.
    -->

    <instance name='default' enabled='false'>
      <!--
	  Wait for network interfaces to be initialized.
      -->
      <dependency name='network'
		  grouping='require_all'
		  restart_on='error'
		  type='service'>
	<service_fmri value='svc:/milestone/network:default'/>
      </dependency>

      <!--
	  Wait for all local filesystems to be mounted.
      -->
      <dependency name='filesystem-local'
		  grouping='require_all'
		  restart_on='none'
		  type='service'>
	<service_fmri
	    value='svc:/system/filesystem/local:default'/>
      </dependency>

      <!--
	  Wait for automounting to be available, as we may be
	  serving data from home directories or other remote
	  filesystems.
      -->
      <dependency name='autofs'
		  grouping='optional_all'
		  restart_on='error'
		  type='service'>
	<service_fmri
	    value='svc:/system/filesystem/autofs:default'/>
      </dependency>

      <exec_method
	  type='method'
	  name='start'
          exec='/usr/sbin/minidlnad -r -f %{config/conf}'
	  timeout_seconds='60' />

      <exec_method
	  type='method'
	  name='stop'
	  exec=':kill'
	  timeout_seconds='15' />

      <exec_method
	  type='method'
	  name='refresh'
	  exec=':kill -HUP'
	  timeout_seconds='15' />

      <property_group name='startd' type='framework'>
	<!-- sub-process core dumps shouldn't restart
	     session -->
	<propval name='ignore_error' type='astring'
		 value='core,signal' />
      </property_group>

      <property_group name='config' type='application'>
	<propval name='conf' type='astring' value='/etc/minidlna.conf' />
      </property_group>

    </instance>

    <stability value='Evolving' />

    <template>
      <common_name>
	<loctext xml:lang='C'>
	  MiniDLNA media server
	</loctext>
      </common_name>
      <documentation>
	<doc_link name='Sourceforge MiniDLNA'
		  uri='https://sourceforge.net/projects/minidlna/' />
      </documentation>
    </template>
  </service>

</service_bundle>
